# debian/simeis.prerm
#!/bin/sh
set -e

if [ "$1" = "remove" ] || [ "$1" = "deconfigure" ]; then
    systemctl stop simeis.service || true
    systemctl disable simeis.service || true
fi

exit 0
