# debian/simeis.postinst
#!/bin/sh
set -e

if ! getent group simeis >/dev/null; then
    addgroup --system simeis
fi

if ! getent passwd simeis >/dev/null; then
    adduser --system --ingroup simeis --no-create-home --home /nonexistent simeis
fi

systemctl daemon-reload || true
systemctl enable simeis.service || true
systemctl start simeis.service || true

exit 0
